# packaging/systemd/hyprwhspr.service
[Unit]
Description=hyprwhspr Voice Dictation Service
After=graphical-session.target pipewire.service
Wants=ydotoold.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/hyprwhspr
Environment=HYPRWHSPR_ROOT=/usr/lib/hyprwhspr
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
