# packaging/systemd/ydotoold.service
[Unit]
Description=ydotool daemon (user)
Documentation=man:ydotoold(1)
After=graphical-session.target
ConditionPathExists=/dev/uinput

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/ydotoold --socket-path=/run/user/%U/ydotool.sock
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
