# config/systemd/ydotoold.service
[Unit]
Description=Generic input server and automation daemon for ydotool (HyprWhspr)
Documentation=man:ydotoold(1)
After=wayland-session@hyprland.desktop.target
Wants=wayland-session@hyprland.desktop.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/ydotoold
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=wayland-session@hyprland.desktop.target
